# tools/decision-profile/templates/com.rayyone.decision-profile-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rayyone.decision-profile-backup</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>HOME_PLACEHOLDER/.claude/decision-profile/daily-backup.sh</string>
	</array>

	<!-- Run daily at 03:00 local time. If the Mac is asleep/off at 03:00,
	     launchd runs it once at the next wake. -->
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<!-- Don't fire on load/login; only on the daily schedule. -->
	<key>RunAtLoad</key>
	<false/>

	<key>StandardOutPath</key>
	<string>HOME_PLACEHOLDER/.ry-decision-profile-backup/launchd.out.log</string>
	<key>StandardErrorPath</key>
	<string>HOME_PLACEHOLDER/.ry-decision-profile-backup/launchd.err.log</string>
</dict>
</plist>
